# home/.display
#!/bin/sh
# Display config.
#
### BEGIN INIT INFO
# Provides:          display
# Required-Start:    $remote_fs $syslog $time
# Required-Stop:     $remote_fs $syslog $time
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: One time execute script to configure displays.
# Description:       Configure displays.
### END INIT INFO

#https://wiki.archlinux.org/index.php/xrandr

# Should change if my config is different.
PRIMARY="eDP-1"
SECONDARY="HDMI-1"

# Necessary since i'm running as root.
export HOME=/home/alex

# For debug purposes only.
echo "Running display config." >> /home/alex/.debug

export DISPLAY=:0

echo $DISPLAY >> /home/alex/.debug

xrandr --display :0 >> /home/alex/.debug

if xrandr | grep "$SECONDARY disconnected"; then
    echo "Disconnected $SECONDARY" >> /home/alex/.debug
    xrandr --output $SECONDARY --off --output $PRIMARY --auto
else
    xrandr --output $SECONDARY --mode 1920x1080 >> /home/alex/.debug
    xrandr --output $PRIMARY --mode 1366x768 >> /home/alex/.debug
    xrandr --output $SECONDARY --left-of $PRIMARY >> /home/alex/.debug
fi
# Here I run the script for change wallpaper. Should not be here!
# /home/alex/.change_wallpaper.sh
